# scripts/template.service.in
[Unit]
Description=@PROJECT_NAME@ service

[Service]
Type=exec
ExecStart=@EXE_PATH@ --config @CMAKE_INSTALL_PREFIX@/@PROJECT_NAME@/share/config.json

[Install]
WantedBy=multi-user.target
