# scripts/template.service.in
[Unit]
Description=@PROJECT_NAME@ service

[Service]
Type=exec
ExecStart=@EXE_PATH@ --config @CMAKE_INSTALL_PREFIX@/share/@PROJECT_NAME@/config.json --enable-prometheus localhost:2212 --telnet-server 2215

[Install]
WantedBy=multi-user.target
